# systemd/pull-backup.timer
[Unit]
Description=Daily pull-backup run

[Timer]
OnCalendar=03:30:00
RandomizedDelaySec=30m
Persistent=true
Unit=pull-backup.service

[Install]
WantedBy=timers.target

